# iot_dashboard_frontend.service
[Unit]
Description=OBSOLETE - React frontend is served by iot_dashboard.service
After=network.target

[Service]
Type=oneshot
User=bokwe
WorkingDirectory=/home/bokwe/IoT_Dashboard/client
# This service is no longer needed; leaving it as a no-op helper.
# Build the frontend manually when deploying:
#   cd /home/bokwe/IoT_Dashboard/client && npm run build
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
